# runit/service/wireplumber/run
#!/bin/sh
exec 2>&1

# Wait for pipewire to be ready
sv check pipewire >/dev/null 2>&1 || exit 1
sleep 1

exec wireplumber
