# runit/service/pulseaudio/run
#!/bin/sh
# PulseAudio user service
exec pulseaudio --daemonize=no --log-target=stderr

