# service/emacs/run
#!/bin/sh
# Start Emacs daemon as this user
exec 2>&1
exec emacs --daemon
